# ip-shuffle.cron
*/3 * * * *  /usr/local/sbin/ip-shuffle
